# home/dot_local/share/dots-ai/dev-companion/systemd/user/dots-ai-dev-companion-worker.service
[Unit]
Description=dots-ai dev companion queue worker (oneshot)
[Service]
Type=oneshot
ExecStart=%h/.local/share/dots-ai/dev-companion/worker.sh
